# recipes-graphics/xorg-session/files/xorg.service
[Unit]
Description=Systemd Xorg Session
Before=step-spu.service
After=systemd-user-sessions.service systemd-logind.service

[Service]
Environment=DISPLAY=:0
TTYPath=/dev/tty1
StandardInput=tty-force
TTYReset=yes
TTYVHangup=yes
TTYVTDisallocate=yes
StandardOutput=journal
StandardError=journal
ExecStart=/usr/bin/Xorg :0 vt1 -keeptty -nolisten tcp -noreset -ac
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
