# recipes-spu/step-spu/files/step-spu.service.in
[Unit]
Description=STEP-SPU Qt App
After=xorg.service network.target
Requires=xorg.service

[Service]
User=$SPU_USER
Group=$SPU_USER
WorkingDirectory=$SPU_WORK_DIR

TTYPath=/dev/tty1
Environment=DISPLAY=:0
Environment=QT_IM_MODULE=qtvirtualkeyboard
Environment=PYTHONDONTWRITEBYTECODE=1
StandardOutput=journal
StandardError=journal
ExecStart=$SPU_APP
Restart=always
RestartSec=1

[Install]
WantedBy=multi-user.target
